# cyberherd.service
[Unit]
Description=Cyber Herd Service

[Service]
ExecStart=/home/sat/bin/herd.sh
