# cyberherd.timer
[Unit]
Description=Restart CyberHerd Service at midnight

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
